# rc.d-init-script
#!/bin/sh

. /etc/rc.subr

name="website"
desc="personal website server"
rcvar="personal_website_enable"
logs="/var/run/website/info.log"
err="/var/run/website/error.log"
command="/usr/local/bin/website"
pidfile="/var/run/website/website.pid"

load_rc_config "website"
run_rc_command "$1"
