# rc.d-init-script
#!/bin/sh

. /etc/rc.subr

name="personal-website"
desc="personal website server"
rcvar="personal_website_enable"

pidfile="/var/run/personal-website/pid"
lockfile="/var/run/personal-website/lock"
command="/home/freebsd/personal-website/start.sh"

load_rc_config "personal-website"
run_rc_command "$1"
